# weather.service
[Unit]
Description=Weather collection daemon for Sense Hat

[Service]
ExecStart=/usr/local/bin/weather

[Install]
WantedBy=multi-user.target